# etc/profile.d/stack-it.sh
export STACKIT_CONF=/etc/stack-it
export STACKIT_ENV=production
